# backend/ankup-core.service
[Unit]
Description=Anklav update system
After=network.target postgresql.service salt-master.service

[Service]
User=ankup_user
Group=ankup_user
WorkingDirectory=/opt/ankup-core
EnvironmentFile=/etc/sysconfig/ankup-core

ExecStartPre=/opt/ankup-core/venv/bin/python manage.py migrate --noinput
ExecStart=/opt/ankup-core/venv/bin/gunicorn --workers 3 --bind 0.0.0.0:8000 core.wsgi:application

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
